# festivald.socket.in
[Socket]
ListenStream=@festivald_socket_path@
SocketUser=@systemd-socket-user@
SocketGroup=@systemd-socket-group@
SocketMode=@systemd-socket-perm@

[Install]
WantedBy=sockets.target
