# festivald.service.in
[Unit]
Description=Festivald speech synthesis service

[Service]
Type=simple
EnvironmentFile=-@festivald_conf_dir@/festivald.conf
ExecStart=@bindir@/festivald
Sockets=festivald.socket
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target

