# festivald.socket.in
[Socket]
ListenStream=@runstatedir@/festivald.socket
SocketUser=@systemd-socket-user@
SocketGroup=@systemd-socket-group@
SocketMode=@systemd-socket-perm@

[Install]
WantedBy=sockets.target
